# services/nm-applet.service
[Unit]
Description=NetworkManager Applet
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/nm-applet --indicator

[Install]
WantedBy=sway-session.target
